# reboot_xapi.sh
#!/bin/bash
### BEGIN INIT INFO
# Provides:          xenserver xapi server
# Required-Start:    from cronta
# Arthur:   Timbaland
# Date:     May-29-2018
# Description:   restart xapi every week
### END INIT INFO
commd = "xapi"
echo $SHELL
#service $commd restart
